# .config/systemd/user/snap-sync.timer
[Unit]
Description=Run snap-sync daily

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
